# cookbooks/xremap/files/.config/autostart/xremap.desktop
[Desktop Entry]
Name=xremap
Exec=systemctl --user start xremap
Terminal=false
Type=Application
